# cron/crontab
# Executes Python Script Every 5 Minutes
@reboot /usr/app/pgaScores.py
*/5 * * * * /usr/app/pgaScores.py
