# cron/crontab
# Executes Python Script Every 5 Minutes
*/5 * * * * /app/pgaScores.py >> /logs/cron.log 2>&1
