# cron/crontab
# Executes Python Script Every 5 Minutes
@reboot /app/pgaScores.py
*/5 * * * * /app/pgaScores.py
